# data/tablet-assist-agent.service
[Unit]
Description=Tablet-mode assistance DBus agent

[Service]
Type=dbus
BusName=tablet.assist.Agent
ExecStart=/usr/sbin/tablet-assist-agent
